# packages/sx05re/emuelec/system.d/ee_shutdown.service
[Unit]
Description=Shutdown EmuELEC Clean
DefaultDependencies=no
Requires=emustation.service retroarch.service
After=emustation.service retroarch.service

[Service]
ExecStart=
ExecStop=/usr/bin/emuelec-utils ee_shutdown
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=emustation.service retroarch.service
